# image/config/image/device/rootfs-overlay/etc/systemd/system/tailscale.service
[Unit]
Description=Tailscale VPN Service
After=network.target systemd-networkd-wait-online@eth0.service overlayfs-check.service systemd-resolved.service
Requires=overlayfs-check.service systemd-networkd-wait-online@eth0.service
Before=hall-sensor-server.service kiosk.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tailscale-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target 