# client/image/device/rootfs-overlay/etc/systemd/system/hall-sensor-server.service
[Unit]
Description=Hall Sensor Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/hall-sensor-server
ExecStart=/usr/bin/node server.js
Restart=always
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target 